# deploy/jobagent-ingest.timer
[Unit]
Description=Frequent job ingestion (runs jobagent-ingest.service every 4h)

[Timer]
OnCalendar=*-*-* 00/4:00:00    # every 4 hours, on the hour
Persistent=true
RandomizedDelaySec=300         # jitter to be polite to sources

[Install]
WantedBy=timers.target
